# etc/profile.d/01-aliases.sh
alias ls='ls -G'
alias ll='ls -lahF'
alias containerfy='cd ~/Tools/containerfy'
alias cwd='realpath .'