# etc/profile.d/05-aliases.sh
# Standard
alias ls='ls -F --color'
alias ll='ls -lah'
alias la='ls -A'
alias l='ls -CF'
alias grep='grep --color=auto'
alias fgrep='fgrep --color=auto'
alias egrep='egrep --color=auto'
# Custom
alias containerfy='cd ~/Tools/containerfy'
alias pwdr='realpath'
alias cdr='cd $(realpath)'
