# cron/root/clean-ds-store.cron
# Only removed .DS_Store from /Volumes folder.
# Hourly (At minute 0)
0 * * * * find /Volumes -name '.DS_Store' -type f -delete >/dev/null 2>&1
# Follow symlinks and recursiveloy remove all occurrences. Not recommended since folder view attributes will be lost.
# 0 * * * * find -L /Volumes -name '.DS_Store' -type f -delete >/dev/null 2>&1
